# sv/loom-provisioner/run
#!/bin/sh
exec 2>&1
exec ruby /opt/loom/provisioner/daemon/provisioner.rb --uri http://localhost:55054 -L info
